# supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/ubuntu/django/all_safe_djangoapp/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
